# dotfiles/.config/autostart/guake.desktop
[Desktop Entry]
Type=Application
Name=Guake
Comment=Start Guake
Exec=guake
Terminal=false
Hidden=false
