# dotfiles/.config/autostart/ibus-autostart.desktop
[Desktop Entry]
Type=Application
Name=iBus daemon
Comment=Start iBus daemon
Exec=ibus-daemon -drx
Icon=ibus
Terminal=false
Hidden=false
